# podman/bbs-pod.service
[Unit]
Description=Daqian BBS Podman Pod (PostgreSQL + App + Nginx)
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=POD_NAME=bbs-pod
EnvironmentFile=-/etc/bbs-pod.env
ExecStart=/usr/bin/podman pod start ${POD_NAME}
ExecStop=/usr/bin/podman pod stop ${POD_NAME}
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
